# config/supervisor/supervisord.ini
[supervisord]
nodaemon=true

[program:nginx]
command=nginx
process_name=%(program_name)s
numprocs=1
directory=/etc/nginx
umask=022
priority=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=1
startretries=1
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=5
serverurl=AUTO
stdout_logfile=/var/log/nginx/access.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/nginx/error.log
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm8
process_name=%(program_name)s
numprocs=1
directory=/etc/php8
umask=022
priority=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=1
startretries=1
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=5
serverurl=AUTO
stdout_logfile=/var/log/php8/access.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/php8/error.log
stderr_logfile_maxbytes=0